# check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/usr/local/bin/check-battery.sh
